# sid/lice_counting/skip_classifier/retraining/crontab
30 15 * * * bash /root/skip_classifier/retraining/run-skip-classifier-retraining.sh >> /var/log/cron.log 2>&1
